# git-build/config/lwsmd.service
[Unit]
Description=Likewise Service Control Manager Service
Wants=network-online.target
After=syslog.target network.target network-online.target haveged.service
Requires=haveged.service

[Service]
Type=forking
ExecStart=PREFIX_DIR/sbin/lwsmd --start-as-daemon --syslog
ExecStartPost=PREFIX_DIR/bin/lwsm autostart
ExecStop=PREFIX_DIR/bin/lwsm shutdown
KillMode=process

[Install]
WantedBy=multi-user.target
